# crontab
* * * * * export SCRAPY=/usr/local/bin/scrapy && /app/runjob.sh >> /var/log/cron.log 2>&1
